# crontab/crontab
PATH=/usr/local/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin

* * * * * bash ~/.crontab/streamdeck.sh
